# systemd/opera-proxy.socket
[Unit]
Description=Opera Proxy Socket
PartOf=opera-proxy.service

[Socket]
ListenStream=127.0.0.1:18080

[Install]
WantedBy=sockets.target
